# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1
#* * * * * laradock /usr/bin/php /var/www/screensavers/artisan schedule:run >> /dev/null 2>&1
#* * * * * laradock /usr/bin/php /var/www/q/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/cartrack/artisan schedule:run >> /dev/null 2>&1
#* * * * * laradock /usr/bin/php /var/www/greeta-laravel/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/myoffer/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/flooringlab/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/macutex/artisan schedule:run >> /dev/null 2>&1
