# laravel-horizon/supervisord.d/laravel-horizon-fpaa.conf
[program:laravel-horizon-fpaa]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/fpaa/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=laradock
stdout_logfile=/var/www/fpaa/storage/logs/horizon.log
