# laravel-horizon/supervisord.d/laravel-horizon-screensavers.conf
[program:laravel-horizon-screensavers]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/screensavers/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=laradock
stdout_logfile=/var/www/screensavers/storage/logs/horizon.log
