# laravel-horizon/supervisord.d/laravel-horizon-phs.conf
[program:laravel-horizon-phs-webapp]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/phs-webapp/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=laradock
stdout_logfile=/var/www/phs-webapp/storage/logs/horizon.log
