# laravel-horizon/supervisord.d/laravel-horizon-sos.conf
[program:laravel-horizon-sos]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/senseofself/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=laradock
stdout_logfile=/var/www/senseofself/storage/logs/horizon.log
