# laravel-horizon/supervisord.d/laravel-horizon-asbestosedge.conf
[program:laravel-horizon-asbestosedge]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php7 /var/www/edgegroup.stealth/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=laradock
stdout_logfile=/var/www/edgegroup.stealth/storage/logs/horizon.log
