# laravel-horizon/supervisord.d/laravel-horizon-kanopi.conf
[program:laravel-horizon-kanopi]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/kanopi/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=laradock
stdout_logfile=/var/www/kanopi/storage/logs/horizon.log
